# deploy/cot-update.timer
[Unit]
Description=COT Data Weekly Update Timer
Requires=cot-update.service

[Timer]
# Every Saturday at 00:00 UTC
OnCalendar=Sat *-*-* 00:00:00
# If missed (VM was off), run on next boot
Persistent=true
# Random delay up to 30 min to avoid thundering herd
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
