# deploy/cot-tgbot.service
[Unit]
Description=COT Analyzer — Telegram Bot & Analytics
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=cftc
Group=cftc
WorkingDirectory=/opt/cftc/backend
ExecStart=/opt/cftc/venv/bin/python tg_bot.py
Restart=always
RestartSec=10

# Environment — set your tokens here or use EnvironmentFile
EnvironmentFile=/opt/cftc/.env
# Alternatively, uncomment and set directly:
# Environment="TG_BOT_TOKEN=your-token-here"
# Environment="TG_CHAT_ID=your-chat-id"
# Environment="ANALYTICS_PORT=8700"

StandardOutput=append:/opt/cftc/backend/logs/tg_bot_systemd.log
StandardError=append:/opt/cftc/backend/logs/tg_bot_systemd.log

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=/opt/cftc
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
